# ops/image-bed-backend.service
[Unit]
Description=Image Bed Backend
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
WorkingDirectory=/root/image-bed/backend
Environment=NODE_ENV=production
ExecStart=/usr/bin/node dist/main
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
