# ops/picvault-backup.timer
[Unit]
Description=Run PicVault Backup Daily

[Timer]
OnCalendar=*-*-* 03:30:00
Persistent=true

[Install]
WantedBy=timers.target
