# ops/picvault-backup.service
[Unit]
Description=PicVault Backup

[Service]
Type=oneshot
WorkingDirectory=/opt/picvault
EnvironmentFile=/opt/picvault/.env
ExecStart=/opt/picvault/scripts/backup.sh
