# ops/picvault-frontend.service
[Unit]
Description=PicVault Frontend
After=network.target picvault-backend.service
Requires=picvault-backend.service

[Service]
Type=simple
WorkingDirectory=/opt/picvault/frontend
Environment=NODE_ENV=production
Environment=HOST=0.0.0.0
Environment=PORT=7899
ExecStart=/usr/bin/node server.mjs
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
